# proxmox-fancontrol.service
[Unit]
Description=Proxmox Fan Control Service
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/root/proxmox-fan-startup.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
